# os-file/files/etc/init.d/autorun
#!/bin/sh /etc/rc.common
# Copyright (C) 2015 OpenWrt.org

START=99

start()
{
	[ -f /root/autorun.sh ] && /root/autorun.sh
}
